# .config/systemd/user/eventd-watcher.path
[Path]
PathModified=%h/.config/eventd

[Install]
WantedBy=default.target
